# with-proxy/services/starman/run
#!/usr/bin/with-contenv /bin/bash

s6-setuidgid www-data /usr/local/bin/run.sh
